# ops/systemd/sanipep-sh-web.service.example
[Unit]
Description=saniPEP SH Web Docker Compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/staff-1/saniPEP-SH-Web
ExecStart=/usr/bin/docker compose up -d --build --remove-orphans
ExecReload=/usr/bin/docker compose up -d --build --remove-orphans
ExecStop=/usr/bin/docker compose down
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
